# fourfor.service
[Unit]
Description="FourFor display"

[Service]
ExecStart=/opt/FourFor/start.sh
Restart=always

[Install]
WantedBy=multi-user.target
