# LidarSlip/95-lidar.rules
KERNEL=="/dev/ttyUSB*", RUN+="/opt/FourFor/LidarSlip/Read-simple.py $KERNEL"
