# LidarSlip/95-lidar.rules
ACTION=="add", KERNEL=="ttyUSB*", ATTRS{idVendor}=="10c4", ATTRS{idProduct}=="ea60", RUN+="/opt/FourFor/LidarSlip/ReadUSB-simple.py %k"
