# LidarSlip/lidar.service
[Unit]
Description="Lidar service"

[Service]
ExecStart=/opt/FourFor/LidarSlip/ReadUSB.py

[Install]
WantedBy=multi-user.target


