# disklavier/midi_over_udp.service
[Unit]
Description="RTP Midi service"

[Service]
ExecStart=/opt/FourFor/disklavier/start.sh
Restart=always

[Install]
WantedBy=multi-user.target


